# src/malamute.service.in
[Unit]
Description=malamute service
After=network.target

[Service]
Type=simple
ExecStart=@prefix@/bin/malamute -f @sysconfdir@/@PACKAGE@/malamute.cfg

[Install]
WantedBy=multi-user.target
